# gitlab-backup.timer
[Unit]
Description=Daily gitlab backup

[Timer]
OnCalendar=*-*-* 2:00:00
AccuracySec=1h
Persistent=true

[Install]
WantedBy=multi-user.target
